# server/kiosk/systemd/engirent-kiosk-browser.service
[Unit]
Description=EngiRent Kiosk Browser (Chromium fullscreen)
After=graphical-session.target engirent-kiosk.service
Wants=engirent-kiosk.service
PartOf=graphical-session.target

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sleep 8
ExecStart=/usr/bin/chromium \
  --noerrdialogs \
  --disable-infobars \
  --disable-session-crashed-bubble \
  --disable-restore-session-state \
  --enable-features=UseOzonePlatform \
  --ozone-platform-hint=auto \
  --kiosk \
  --app=http://localhost:8080 \
  --window-position=0,0 \
  --start-fullscreen \
  --check-for-update-interval=31536000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=engirent-browser

[Install]
WantedBy=graphical-session.target
